# systemd/cyan-dash.service
[Unit]
Description=Lume Dashboard
After=network.target

[Service]
Type=simple
# Update these paths to match your installation
WorkingDirectory=%h/Lume
EnvironmentFile=%h/Lume/.env
ExecStart=/usr/bin/node server/server.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
