# systemd/cyan-dash.service
[Unit]
Description=Cyan Dashboard
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/workspace/projects/cyan-dash/server
EnvironmentFile=%h/workspace/projects/cyan-dash/.env
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
